# sentinel-agent/sentinel-agent.service
[Unit]
Description=Sentinel Agent — remote telemetry and patch relay
After=network.target

[Service]
Type=simple
User=sentinel-agent
WorkingDirectory=/opt/sentinel-agent
EnvironmentFile=/etc/sentinel-agent/env
ExecStart=/opt/sentinel-agent/venv/bin/python /opt/sentinel-agent/agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-agent

[Install]
WantedBy=multi-user.target
